# .config/systemd/user/scroll-debounce.service
[Unit]
Description=Scroll debounce - drop phantom wheel ticks from a worn mouse encoder
After=graphical-session.target

[Service]
ExecStart=%h/.local/bin/scroll-debounce
Restart=always
RestartSec=2
# Debounce window in ms. Raise if phantom ticks still slip through; lower if
# fast intentional scroll reversals get eaten.
Environment=SD_WINDOW_MS=50
# Uncomment to log every dropped/kept reversal to the journal for tuning:
#   journalctl --user -u scroll-debounce -f
# Environment=SD_DEBUG=1

[Install]
WantedBy=default.target
