# .config/systemd/user/swayosd-server.service.d/override.conf
# Force GTK4's cairo renderer for swayosd-server.
#
# Why: GTK4's default GL/NGL renderer aliases rounded-corner clips — the OSD
# pill's `border-radius` shows visible stair-step "ladders" on the curve,
# especially over high-contrast backgrounds and on this host's fractional
# monitor scale (1.3). The cairo renderer anti-aliases the corner cleanly.
# Verified: GL renderer = jagged arc, cairo = smooth arc (same CSS).
#
# Pairs with the frosted-glass swayosd/style.css (rounded translucent pill).
# Drop-in overrides the omarchy-provided swayosd-server.service without editing it.
# Apply: ./apply.sh && systemctl --user daemon-reload && omarchy restart swayosd
[Service]
Environment=GSK_RENDERER=cairo
